# test/rootfs/etc/s6-overlay/s6-rc.d/otelcol/run
#!/bin/ash
/usr/bin/otelcol --config /etc/otelcol/config.yaml
