# test/s6/app/run
#!/bin/ash
/usr/bin/php /usr/local/src/bin/hyperf.php start
